# ic_recorder.service
[Service]
ExecStart=/home/ben/projects/ic_recorder/ic_recorder

[Unit]
Description="IC Recorder Service"

[Install]
WantedBy=multi-user.target
