# ic_recorder.path
[Path]
PathExists=-media-ben-IC\x20RECORDER

[Unit]
Description="IC Recorder Path"
PartOf=ic_recorder.service
