# ic_recorder.path
[Unit]
Description="IC Recorder Path"

[Path]
PathExists=/media/ben/IC RECORDER

[Install]
WantedBy=multi-user.target
