# ic_recorder.service
[Unit]
Description="IC Recorder Service"

[Service]
ExecStart=/usr/bin/ic_recorder

[Install]
WantedBy=multi-user.target
